# root/etc/services.d/nginx/run
#!/bin/execlineb
nginx -g "daemon off;pid /run/nginx.pid;" -c /etc/nginx/nginx.conf